# systemd/sirens-discord-ops-update.timer
[Unit]
Description=sirens-discord-ops auto-update poll (every 5min)

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=30s
Unit=sirens-discord-ops-update.service

[Install]
WantedBy=timers.target
