# deploy/sirens-discord-ops.service
[Unit]
Description=Sirens Discord ops bot (eco/factorio admin control panel)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=kai
Group=kai
EnvironmentFile=/etc/sirens-discord-ops.env
ExecStart=/usr/local/bin/sirens-discord-ops
Restart=on-failure
RestartSec=5s

# The bot shells out to `coily gaming eco <verb>`, so it needs the same
# environment as a coily-on-kai-server invocation. coily handles its own
# privilege escalation through sudo where required.
NoNewPrivileges=false
PrivateTmp=true
ProtectHome=read-only
ProtectSystem=full

[Install]
WantedBy=multi-user.target
